# scripts/systemd/system/ogscope.service.d/ogscope-hardware-plane.conf
[Unit]
# 为后续拆分独立进程预留依赖 / Reserve dependency for future split process
Wants=ogscope-hardware-plane.service
After=ogscope-hardware-plane.service

[Service]
# 共享硬件平面控制开关 / Shared hardware plane switch
Environment=OGSCOPE_HARDWARE_PLANE_ENABLED=true
