# scripts/systemd/system/ogscope.service.d/ogscope-hardware-plane.conf
# 共享硬件平面：进程内实现，仅注入环境变量（无独立 unit 依赖）
# Shared hardware plane runs in-process; environment only (no separate unit dependency)
#
# 安装：由 scripts/install.sh 与 scripts/board-update.sh 同步到
# /etc/systemd/system/ogscope.service.d/ogscope-hardware-plane.conf
# Installed by install.sh and board-update.sh to ogscope.service.d/

[Service]
Environment=OGSCOPE_HARDWARE_PLANE_ENABLED=true
