# scripts/systemd/system/ogscope.service.d/ogscope-low-ram.conf
# Raspberry Pi Zero 2W（约 416MiB RAM）额外节流 / Extra throttling for ~416MiB boards
# 默认采集与提星长边已由 config 全局设为 1280×720 / Base 1280×720 is global in config
# 关闭自动安装：OGSCOPE_SKIP_LOW_RAM_DEFAULTS=1 / Disable: OGSCOPE_SKIP_LOW_RAM_DEFAULTS=1

[Service]
Environment=OGSCOPE_SOLVER_MAX_STARS_HARD_CAP=40
Environment=OGSCOPE_SOLVER_MAX_IMAGE_SIDE_HARD_CAP=1280
Environment=OGSCOPE_SHARED_PREVIEW_FPS=5
Environment=OGSCOPE_PREVIEW_JPEG_QUALITY=65
Environment=OGSCOPE_STREAM_MAX_MJPEG_CLIENTS=1
