# scripts/systemd/system/ogscope.service.d/ogscope-memory.example.conf
# 示例：限制 ogscope 单元 cgroup 内存（需复制为 memory-limit.conf 或安装时设 OGSCOPE_SYSTEMD_MEMORY_MAX）
# Example: cap ogscope unit memory (copy to memory-limit.conf or set OGSCOPE_SYSTEMD_MEMORY_MAX at install)
#
# 注意：过小会在 cgroup 内 SIGKILL；过大仍可能整机 OOM。常与 OGSCOPE_SOLVER_*_HARD_CAP 同用。
# Note: too low kills the service early; too high may still OOM the system. Pair with solver hard caps.

[Service]
# MemoryMax=380M
